# templates/default/init/jira.erb
#!/bin/sh -e
# JIRA startup script
# chkconfig: 2345 90 10
# description: JIRA

USER=<%= node[:jira][:user] %>
BASE=<%= node[:jira][:app_dir] %>
export JAVA_HOME=<%= node[:java][:home] %>

case "$1" in
  start)
    echo "Starting jira"
    /bin/su -m $USER -c "cd $BASE/logs && $BASE/bin/start-jira.sh &> /dev/null"
    ;;
  
  stop)
    echo "Stopping jira"
    /bin/su -m $USER -c "$BASE/bin/stop-jira.sh &> /dev/null"
    echo "Jira Stopped successfully"
    ;;
  
  restart)
    $0 stop
    sleep 5
    $0 start
    ;;
  
  *)
    echo "Usage: /etc/init.d/jira {start|restart|stop}"
    exit 1
    ;;
esac
exit 0
